# pi-sensor-api.user.service
[Unit]
Description=Pi DHT11 Sensor API
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pi-sensor-api
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/pi-sensor-api.env
WorkingDirectory=/opt/pi-sensor-api

[Install]
WantedBy=default.target
